# systemd/stp-display.service
[Unit]
Description=Sun-Moon Transit Predictor — e-paper display client
Documentation=https://github.com/joergs-git/sun-moon-transit-predictor
# Start after the predictor so /api/config and /api/state are already up when
# the panel first polls. Wants (not Requires) so the display still runs and
# shows "SERVER OFFLINE" if the predictor is down or points at a remote host.
After=network-online.target stp.service
Wants=network-online.target stp.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__/display
ExecStart=/usr/bin/python3 __INSTALL_DIR__/display/epaper_client.py
Restart=on-failure
RestartSec=5
# Kill a wedged panel quickly so `systemctl restart` can't hang on a hung
# ReadBusy and the GPIO is released before the new instance starts (v0.47.3).
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

# The single bootstrap value: where to read the live `display` config block
# from (everything else is configured in the browser Settings panel). Point
# this at a remote host only if the predictor's web UI itself lives elsewhere;
# to merely pull DATA from another Pi, leave this as localhost and set the
# "Data source URL" field in Settings instead.
Environment=STP_CONFIG_URL=http://127.0.0.1:8081
# For older 4.2" boards that show nothing on first run, switch to: epd4in2
Environment=STP_EPD_DRIVER=epd4in2_V2
# Pi 5: use the lgpio pin factory directly (RPi.GPIO/pigpio don't work on the
# Pi 5's RP1 GPIO). Skips gpiozero's noisy fallback chain.
Environment=GPIOZERO_PIN_FACTORY=lgpio

# Hardware access: driving the panel needs /dev/gpiomem*, /dev/gpiochip* and
# /dev/spidev*, and lgpio writes a small .lgd-nfy* notify file in the working
# directory. So this service deliberately does NOT use ProtectHome/ProtectSystem
# /DeviceAllow — those blocked gpiomem and the lgpio notify file and broke every
# pin factory. It stays a non-privileged user in the spi+gpio groups (the
# installer adds them); keep the two lightweight protections that don't impede
# the panel.
NoNewPrivileges=true
PrivateTmp=true
SupplementaryGroups=spi gpio

[Install]
WantedBy=multi-user.target
